# files/podman/podman-startup.service
[Unit]
Description=Start Podman containers in order after system is ready
After=network-online.target netavark-dhcp-proxy.socket local-fs.target nss-lookup.target multi-user.target zfs-mount.service
Wants=network-online.target netavark-dhcp-proxy.socket

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=no
TimeoutStartSec=600
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=podman-startup

# Wait for system to stabilize
ExecStartPre=/bin/sleep 30

# Phase 1: Database containers (parallel start)
ExecStart=-/bin/bash -c '\
echo "=== Phase 1: Starting databases ==="; \
podman start $(podman ps -a --format "{{.Names}}" | grep -E "(-postgres|-mariadb|-valkey|-influxdb|-mongo|-elasticsearch)$"); \
echo "Waiting for databases to initialize..."; \
sleep 10'

# Phase 2: Core infrastructure (sequential start)
ExecStart=-/bin/bash -c '\
echo "=== Phase 2: Starting core services ==="; \
podman start authelia; \
sleep 30; \
podman start traefik; \
sleep 5'

# Phase 3: All remaining containers (batched start)
ExecStart=-/bin/bash -c '\
echo "=== Phase 3: Starting remaining containers ==="; \
podman ps -a --filter "status=exited" --filter "status=created" --format "{{.Names}}" | xargs -n 10 -r podman start; \
sleep 5; \
echo "=== All phases complete ==="'

[Install]
WantedBy=multi-user.target
